# src/JOHNNYbeGOOD.Home.Api/johnnybegood.service
[Unit]
Description=JOHNNYbeGOOD Feeder

[Service]
Type=notify
ExecStart=/usr/local/bin/johnnybegood/JOHNNYbeGOOD.Home.Api

[Install]
WantedBy=multi-user.target